# .config/systemd/user/udiskie.service
[Unit]
Description=Automount drives daemon

[Service]
ExecStart=udiskie --no-notify

[Install]
WantedBy=default.target
